# supervisor.conf
[supervisord]
nodaemon=true

[program:node-superv]
directory=/app
command=npm start
autostart=true
autorestart=true
redirect_stderr=true
environment=NODE_ENV="prod"